# docker-grafana/supervisord.conf.d/supervisord.nginx.conf
[program:nginx]
priority=100
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
command=/usr/sbin/nginx